# ops/systemd/smartbarber-backup.service
[Unit]
Description=Smart Barber Shops backup job
Wants=docker.service
After=docker.service

[Service]
Type=oneshot
User=barberadmin
Group=barberadmin
WorkingDirectory=/opt/smartbarber/app
Environment=BACKUP_RETENTION_DAYS=14
ExecStart=/opt/smartbarber/app/backup.sh

